# ragd/scripts/ragd.service
[Unit]
Description=RAGD - Blackmark Dominion RAG Daemon
After=network.target
Documentation=https://github.com/MatinDeevv/BlackmarkDominion

[Service]
Type=simple
ExecStart=/usr/local/bin/ragd --config %h/.ragd/config.json --daemon
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=HOME=%h
Environment=RAGD_ENDPOINT=http://localhost:7474

[Install]
WantedBy=default.target
